# systemd/system/pen-locker.service
[Unit]
Description=Run pen-locker recv

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/pen-locker no-file --cmd recv