# systemd/system/pen-locker.path
[Unit]
Description=Run pen-locker recv

[Path]
PathModified=/tmp/pen-locker.path

[Install]
WantedBy=multi-user.target
